# .config/systemd/user/udev-keycool-keyboard.service
[Unit]
Description=Keycool keyboard

[Service]
ExecStart=kmonad %h/.config/kmonad/config_keycool.kbd
Type=oneshot
RemainAfterExit=no
