# udevrules/keyboard.rules
SUBSYSTEMS=="usb", DRIVERS=="usb", RUN+="/home/magnickolas/scripts/setup_keyboard.sh"
