# crontab/crontab
* * * * * ~/scripts/battery_notification.sh >/dev/null 2>&1
