# crontab/crontab
* * * * * ~/.local/bin/x battery-notify
